# juno/systemd/imap-backup.service
[Unit]
Description=IMAP Backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mail-backup.sh